# scripts/weave/weave-network.target
[Unit]
Description=Weave Network Setup Complete
Documentation=man:systemd.special(7)
RefuseManualStart=no
[Install]
WantedBy=multi-user.target